# .config/systemd/user/compton.service
[Unit]
Description=Compton

[Service]
ExecStart=/usr/bin/compton --backend xrender

[Install]
WantedBy=default.target
