# .config/systemd/user/dbus.socket
[Unit]
Description=D-Bus User Bus Socket

[Socket]
ListenStream=%t/bus
