# .config/systemd/user/dbus.socket
[Unit]
Description=D-Bus User Bus Socket

[Socket]
ListenStream=%t/dbus/user_bus_socket
#ListenStream=%t/bus
