# .config/systemd/user/default.target
[Unit]
Description=Systemd user default target
After=xorg.target

[Install]
Alias=default.target
