# .config/systemd/user/dwmstatus.service
[Unit]
Description=Dwm Status Line

[Service]
ExecStart=/usr/bin/dwmstatus
Restart=no

[Install]
WantedBy=default.target
