# .config/systemd/user/dwmstatus@.service
[Unit]
Description=Dwm Status Line on %i
BindsTo=sys-subsystem-net-devices-%i.device
After=network.target

[Service]
Environment=MAILDIR=%h/mail
ExecStart=/usr/bin/dwmstatus -i %i

[Install]
WantedBy=default.target
