# .config/systemd/user/firefox.service
[Unit]
Description=Firefox

[Service]
ExecStart=/usr/bin/firefox

[Install]
WantedBy=default.target
