# .config/systemd/user/mpd.socket
[Unit]
Description=Music Player Daemon Socket

[Socket]
ListenStream=%t/mpd/socket

[Install]
WantedBy=sockets.target
