# .config/systemd/user/mpd.socket
[Unit]
Description=Music Player Daemon Socket

[Socket]
ListenStream=6600

[Install]
WantedBy=sockets.target
