# .config/systemd/user/udiskie.service
[Unit]
Description=UDiskie Automounter

[Service]
ExecStart=/usr/bin/udiskie -N -F

[Install]
WantedBy=default.target
