# .config/systemd/user/mpd.service.d/extra.conf
[Unit]
Documentation=man:mpd(1)

[Service]
Environment=PULSE_PROP=application.icon_name=sound

[Install]
Also=mpd.socket
